# satori/images/riemann/app/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
childlogdir=/var/log

[unix_http_server]
file=/var/run//supervisor.sock
chmod=0770
chown=root:root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run//supervisor.sock

[program:riemann]
command = /app/bootstrap.sh
user = root
autostart = false
autorestart = true

[program:reloader]
command = /app/riemann-reloader
user = root
autostart = true
autorestart = true
